# crontab
*/2 * * * * /home/emersonmachado/server_monitoring.sh /home/emersonmachado